# etc/supervisor.d/run.conf
[program:clamspringservice]
command=java -jar app.jar
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0